# docker/images/elasticsearch-6.1.3/ROOT/etc/services.d/discovery/run
#!/usr/bin/with-contenv /bin/sh

set -e 

echo "Starting background discovery..."

exec /srv/container-utils/es-bg-discovery.sh